# 2024/quals/pwn-fasterbox/challenge/chal.target
[Unit]
Description=Chal target
Requires=sysinit.target
Wants=chal.service
After=sysinit.target
AllowIsolate=yes
