# 2024/quals/pwn-fasterbox/challenge/chal.service
[Unit]
Description=Chal service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/chal.sh
ExecStartPost=/usr/bin/systemctl poweroff
